# infra/systemd/iaase-discover.timer
[Unit]
Description=IAASE Facebook Marketplace - Discovery (Timer)

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
RandomizedDelaySec=30s
Persistent=true
Unit=iaase-discover.service

[Install]
WantedBy=timers.target
