# infra/systemd/iaase-monitor.timer
[Unit]
Description=IAASE Facebook Marketplace - Monitor (Timer)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=2min
Persistent=true
Unit=iaase-monitor.service

[Install]
WantedBy=timers.target
